# rootfs/.init/services/run/init_root.conf
[program:test]
command = /build/bin/init_root
startsecs = 0
autorestart = false
startretries = 1