# rootfs/build/services/coder.conf
[program:coder]

command = "/build/bin/start_coder"
autostart = true
autorestart = true
