# rootfs/build/services/coder.conf
[program:coder]

command = "/build/bin/start_code"
autostart = true
autorestart = true
