# rootfs/build/services/jupyter.conf
[program:jupyter]

command = "/build/bin/start_jupyter"
autostart = true
autorestart = true