# rootfs/build/services/sshd.conf
[program:sshd]

command = "/build/bin/start_sshd"
autostart = true
autorestart = true
