# rootfs/build/services/vscode.conf
[program:vscode]

command = "/build/bin/start_vscode"
autostart = true
autorestart = true
